# rootfs/config/supervisor/2-vnctiger.conf
[program:VNCTiger]
command=bash -c 'vncserver -depth 32 -geometry $VNC_RESOLUTION :0'
startsecs=0
autorestart=unexpected
stderr_logfile=/var/log/VNCTiger.stderr.log
priority=2
