# rootfs/config/supervisor/3-vncproxy.conf
[program:TVNCproxy]
command=bash -c '/opt/noVNC/utils/novnc_proxy --vnc 0.0.0.0:5900 --listen 80'
startsecs=0
autorestart=unexpected
stderr_logfile=/var/log/TVNCproxy.stderr.log
priority=3