# rootfs/config/supervisor/6-noVNC.conf
[program:vncproxy]
command=bash -c '/opt/noVNC/utils/novnc_proxy --vnc 0.0.0.0:5900 --listen 90'
startsecs=0
autorestart=unexpected
stderr_logfile=/var/log/novncproxy.stderr.log
priority=4